# contrib/dumpstore.rc
#!/bin/sh
#
# PROVIDE: dumpstore
# REQUIRE: LOGIN FILESYSTEMS zfs
# KEYWORD: shutdown
#
# Add to /etc/rc.conf to enable:
#   dumpstore_enable="YES"
#
# Optional rc.conf knobs (with defaults shown):
#   dumpstore_addr=":8080"
#   dumpstore_dir="/usr/local/lib/dumpstore"
#   dumpstore_config="/etc/dumpstore/dumpstore.conf"
#   dumpstore_logfile="/var/log/dumpstore.log"

. /etc/rc.subr

name="dumpstore"
rcvar="${name}_enable"
desc="ZFS web UI"

load_rc_config ${name}

: ${dumpstore_enable:="NO"}
: ${dumpstore_addr:=":8080"}
: ${dumpstore_dir:="/usr/local/lib/dumpstore"}
: ${dumpstore_config:="/etc/dumpstore/dumpstore.conf"}
: ${dumpstore_logfile:="/var/log/dumpstore.log"}

# procname is the real binary; command is the daemon(8) supervisor.
# rc.subr uses procname for status/stop checks.
procname="${dumpstore_dir}/dumpstore"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-P ${pidfile} -r -o ${dumpstore_logfile} -m 3 \
    ${procname} \
    -addr ${dumpstore_addr} \
    -dir ${dumpstore_dir} \
    -config ${dumpstore_config}"

run_rc_command "$1"
