# contrib/dumpstore.rc
#!/bin/sh
#
# PROVIDE: dumpstore
# REQUIRE: LOGIN FILESYSTEMS zfs
# KEYWORD: shutdown
#
# Add to /etc/rc.conf to enable:
#   dumpstore_enable="YES"
#
# Optional rc.conf knobs (with defaults shown):
#   dumpstore_enable (bool):   Set to YES to enable dumpstore.
#   dumpstore_addr (str):      Listen address. Default: ":8080"
#   dumpstore_dir (path):      Base directory. Default: /usr/local/lib/dumpstore
#   dumpstore_config (path):   Config file. Default: /etc/dumpstore/dumpstore.conf
#   dumpstore_logfile (path):  Log file. Default: /var/log/dumpstore.log

. /etc/rc.subr

name="dumpstore"
rcvar="${name}_enable"
desc="ZFS web UI"

load_rc_config ${name}

: ${dumpstore_enable:="NO"}
: ${dumpstore_addr:=":8080"}
: ${dumpstore_dir:="/usr/local/lib/dumpstore"}
: ${dumpstore_config:="/etc/dumpstore/dumpstore.conf"}
: ${dumpstore_logfile:="/var/log/dumpstore.log"}

pidfile="/var/run/${name}.pid"
procname="${dumpstore_dir}/dumpstore"
command="/usr/sbin/daemon"
command_args="-p ${pidfile} -o ${dumpstore_logfile} -m 3 -S -T ${name} \
    ${procname} \
    -addr ${dumpstore_addr} \
    -dir ${dumpstore_dir} \
    -config ${dumpstore_config}"

export PATH="${PATH}:/usr/local/bin:/usr/local/sbin"

start_precmd="${name}_prestart"

dumpstore_prestart()
{
    if [ ! -e ${pidfile} ]; then
        install -m 0600 /dev/null ${pidfile}
    fi
}

run_rc_command "$1"
