# contrib/dumpstore.rc
#!/bin/sh
#
# PROVIDE: dumpstore
# REQUIRE: LOGIN FILESYSTEMS zfs
# KEYWORD: shutdown
#
# Add to /etc/rc.conf to enable:
#   dumpstore_enable="YES"
#
# Optional rc.conf knobs (with defaults shown):
#   dumpstore_addr=":8080"
#   dumpstore_dir="/usr/local/lib/dumpstore"
#   dumpstore_config="/etc/dumpstore/dumpstore.conf"
#   dumpstore_logfile="/var/log/dumpstore.log"

. /etc/rc.subr

name="dumpstore"
rcvar="${name}_enable"
desc="ZFS web UI"

load_rc_config ${name}

: ${dumpstore_enable:="NO"}
: ${dumpstore_addr:=":8080"}
: ${dumpstore_dir:="/usr/local/lib/dumpstore"}
: ${dumpstore_config:="/etc/dumpstore/dumpstore.conf"}
: ${dumpstore_logfile:="/var/log/dumpstore.log"}

_bin="${dumpstore_dir}/dumpstore"
_supervisor_pid="/var/run/${name}-supervisor.pid"
_child_pid="/var/run/${name}.pid"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
restart_cmd="${name}_restart"
extra_commands="status"

dumpstore_start() {
    echo "Starting ${name}."
    /usr/sbin/daemon \
        -P ${_supervisor_pid} \
        -p ${_child_pid} \
        -r \
        -o ${dumpstore_logfile} \
        -m 3 \
        ${_bin} \
        -addr ${dumpstore_addr} \
        -dir ${dumpstore_dir} \
        -config ${dumpstore_config}
}

dumpstore_stop() {
    echo "Stopping ${name}."
    if [ -f "${_supervisor_pid}" ]; then
        kill "$(cat ${_supervisor_pid})" 2>/dev/null || true
        rm -f "${_supervisor_pid}"
    fi
    if [ -f "${_child_pid}" ]; then
        kill "$(cat ${_child_pid})" 2>/dev/null || true
        rm -f "${_child_pid}"
    fi
}

dumpstore_status() {
    if [ -f "${_supervisor_pid}" ] && kill -0 "$(cat ${_supervisor_pid})" 2>/dev/null; then
        echo "${name} is running as pid $(cat ${_child_pid} 2>/dev/null)."
    else
        echo "${name} is not running."
        return 1
    fi
}

dumpstore_restart() {
    dumpstore_stop
    sleep 1
    dumpstore_start
}

run_rc_command "$1"
